# configs/nftables/archmailbox-nftables.service
[Unit]
Description=Archmailbox dkim proxy filter
Documentation=man:nft(8)
Wants=network-pre.target
Before=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nft -f /etc/archmailbox-nftables.conf
ExecReload=/usr/bin/nft delete table ip dkim_filter ';' delete table ip6 dkim_filter ';' include '"/etc/archmailbox-nftables.conf"'
ExecStop=/usr/bin/nft delete table ip dkim_filter ';' delete table ip6 dkim_filter
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
